# hawkeye2ha/hawkeye2ha/rootfs/etc/services.d/hawkeye2ha/finish
#!/usr/bin/with-contenv bashio
if [[ "${1}" != "0" ]]; then
    bashio::log.warning "hawkeye2ha exited with code ${1}"
fi
