# hawkeye2ha/hawkeye2ha/rootfs/etc/services.d/hawkeye2ha/run
#!/usr/bin/with-contenv bashio
exec python3 -m hawkeye2ha
